# gcp-gpu-node/cudajulia.path
[Unit]
Description = Install Julia CUDA components when /usr/local/bin/julia is available

[Path]
PathExists = /opt/julia/julia-1.1.0/bin/julia
