# gcp-gpu-node/cudapython.path
[Unit]
Description = Install Python CUDA components when conda is available

[Path]
PathExists = /opt/anaconda3/bin/conda
